# systemd/hermes-github-pr-kanban-bridge.timer
[Unit]
Description=Run Hermes GitHub PR activity to Kanban bridge every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=hermes-github-pr-kanban-bridge.service

[Install]
WantedBy=timers.target
